# rclone-gdrive/root/etc/services.d/rclone_mount/run
#!/usr/bin/with-contenv sh

export LOGFILE="/logs/rclone_mount.log"

exec s6-setuidgid abc /usr/bin/rclone mount rclone-gdrive-union: ${LOCAL_MOUNT} \
--log-level DEBUG --log-file=$LOGFILE \
--cache-dir /data/rclone/cache \
--config /rclone.conf ${RCLONE_MOUNT_OPTIONS}