# rclone-gdrive/root/etc/services.d/rclone_mount/run
#!/usr/bin/with-contenv sh

export LOGFILE="/logs/rclone_mount.log"

s6-setuidgid abc mkdir -p ${REMOTE_MOUNT}

exec s6-setuidgid abc /usr/bin/rclone mount ${REMOTE_PATH} ${REMOTE_MOUNT} \
--log-level INFO --log-file=$LOGFILE --rc \
--config $CONFIG_FILE $RCLONE_MOUNT_OPTIONS