# rclone-gdrive/root/etc/services.d/rclone_mount/run
#!/usr/bin/with-contenv sh

function unmount_all {
    fusermount -uz ${LOCAL_MOUNT} || :
    fusermount -uz ${REMOTE_MOUNT} || :
}
unmount_all 2> /dev/null

mergerfs ${MERGERFS_OPTIONS} ${UPLOAD_PATH}:${REMOTE_MOUNT}  ${LOCAL_MOUNT}

exec /usr/bin/rclone mount --uid $PUID --gid $PGID ${REMOTE_PATH} ${REMOTE_MOUNT} \
    --cache-dir /data/rclone/cache ${RCLONE_MOUNT_OPTIONS}